# 80-beanstalkd.conf
[program:beanstalkd]
command=/usr/bin/beanstalkd -n -u nobody
autostart=true
autorestart=true
