# 90-queue-worker.conf
[program:queue-worker]
command=/usr/bin/php artisan queue:listen
autostart=true
autorestart=true
directory=/home/vagrant/vhosts/gonorthwebsites/
